# etc/systemd/system/getty@tty1.service.d/canute.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --nonewline --issue-file /brltty/current/etc/issue --login-options '-p -- \\u' --noclear %I $TERM
