# etc/systemd/system/brltty@.service.d/wrapper.conf
[Service]
Environment=BRLTTY_EXECUTABLE_PATH=/brltty/current/bin/brltty

