# etc/profile.d/brltty.sh
export PATH="/brltty/current/bin:${PATH}"
