# etc/profile.d/zzz-canute-startx-if-on-tty1.sh
if [ `tty` = "/dev/tty1" ]
then
   exec canute-startx -c
fi
