# Autostart/Systemd/brltty@.service.in
###############################################################################
# BRLTTY - A background process providing access to the console screen (when in
#          text mode) for a blind person using a refreshable braille display.
#
# Copyright (C) 1995-2020 by The BRLTTY Developers.
#
# BRLTTY comes with ABSOLUTELY NO WARRANTY.
#
# This is free software, placed under the terms of the
# GNU Lesser General Public License, as published by the Free Software
# Foundation; either version 2.1 of the License, or (at your option) any
# later version. Please see the file LICENSE-LGPL for details.
#
# Web Page: http://brltty.app/
#
# This software is maintained by Dave Mielke <dave@mielke.cc>.
###############################################################################

[Unit]

Description=BRLTTY Instance: %I
Documentation=man:brltty(1)
Documentation=@PACKAGE_URL@

ConditionPathExists=%I
RequiresMountsFor=@UPDATABLE_DIRECTORY@
RequiresMountsFor=@api_socket_directory@

Wants=systemd-udev-settle.service
After=systemd-udev-settle.service

StopWhenUnneeded=true
RefuseManualStart=true
RefuseManualStop=true


[Service]

ExecStartPre=-/sbin/modprobe pcspkr
ExecStartPre=-/sbin/modprobe uinput

Environment="BRLTTY_SYSTEMD_INSTANCE=%I"
ExecStart=/usr/lib/brltty/systemd-wrapper

Type=@SYSTEMD_SERVICE_TYPE@
Environment="BRLTTY_EXECUTABLE_ARGUMENTS=@SYSTEMD_SERVICE_ARGUMENTS@"

PIDFile=%t/brltty/brltty-%i.pid
Environment="BRLTTY_PID_FILE=%t/brltty/brltty-%i.pid"

TimeoutStartSec=5
TimeoutStopSec=10

Restart=on-failure
RestartSec=30

Nice=-10
OOMScoreAdjust=-900

#######################################################
# keep these synchronized with sysusers.d/brltty.conf #
#######################################################

# for reading screen content (/dev/vcs*)
SupplementaryGroups=tty

# for virtual console access (/dev/tty<n>)
SupplementaryGroups=tty

# for serial I/O (/dev/ttyS<n>)
SupplementaryGroups=dialout

# for USB I/o (/dev/bus/usb/*/*)
SupplementaryGroups=root

# for sound via the ALSA framework
SupplementaryGroups=audio

# for sound via the Pulse Audio daemon
SupplementaryGroups=pulse-access

# for keyboard detection (/dev/input/*)
SupplementaryGroups=input

# for virtual keyboard creation (/dev/uinput)
SupplementaryGroups=root

