# Autostart/Systemd/brltty@.service.in
[Unit]

Description=BRLTTY Instance: %I
Documentation=man:brltty(1)
Documentation=@PACKAGE_URL@

DefaultDependencies=no
StopWhenUnneeded=true

RequiresMountsFor=@UPDATABLE_DIRECTORY@


[Service]

Environment="BRLTTY_SYSTEMD_INSTANCE=%I"
ExecStart=/sbin/brltty-systemd-wrapper

Type=@SYSTEMD_SERVICE_TYPE@
Environment="BRLTTY_EXECUTABLE_ARGUMENTS=@SYSTEMD_SERVICE_ARGUMENTS@"

PIDFile=%t/brltty-%i.pid
Environment=BRLTTY_PID_FILE=%t/brltty-%i.pid

TimeoutStartSec=5
TimeoutStopSec=10

Restart=always
RestartSec=30

Nice=-10
OOMScoreAdjust=-900

ProtectHome=read-only
ProtectSystem=full


