# Autostart/Systemd/brltty@.service.in
[Unit]

Description=Braille Device: %I
Documentation=man:brltty(1)
Documentation=@PACKAGE_URL@

DefaultDependencies=no
StopWhenUnneeded=true

RequiresMountsFor=@UPDATABLE_DIRECTORY@


[Service]

Type=notify
ExecStart=/bin/sh -c 'device="%I" /bin/brltty -n -E -P /run/brltty-%I.pid -d "${device#DEFAULT}"'

TimeoutStartSec=5
TimeoutStopSec=10

Restart=always
RestartSec=30

Nice=-10
OOMScoreAdjust=-900

ProtectHome=read-only
ProtectSystem=full
SystemCallArchitectures=native


